# autostart/set-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=set-wallpaper
NoDisplay=true
Exec=${libexecdir}/nsls2/set-wallpaper
Terminal=false