# autostart/set-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=set-wallpaper
NoDisplay=true
Exec=
Terminal=false